# assets/tool-seeds/workflow/com.agentosai.content-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentosai.content-scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/hoony/projects/agent-os-ai/assets/tool-seeds/workflow/scheduler.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/hoony/projects/agent-os-ai</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/tmp/agentosai-content-scheduler.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/agentosai-content-scheduler.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- claude CLI 위치 보강 (nvm·brew 양쪽 커버) -->
        <key>PATH</key>
        <string>/Users/hoony/.nvm/versions/node/v22.22.2/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <!-- 기본 회차 주기 (시간). .env 의 ROUTINE_INTERVAL_HOURS 가 있으면 그게 우선. -->
        <key>ROUTINE_INTERVAL_HOURS</key>
        <string>4</string>
    </dict>
</dict>
</plist>
